# magento2/tools/docker/supervisor/conf.d/apps.conf
[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
autostart=true
autorestart=true
user=root

[program:sshd]
command = /bin/bash -c "env > /etc/environment && /usr/sbin/sshd -D"
autorestart=true

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true